# builtin/tool/data/graphite/conf/common/supervisord.conf
[supervisord]
nodaemon = true

[program:carbon-relay-1]
command = /usr/bin/carbon-relay --config=/etc/carbon/carbon.conf --instance=1 --debug --pidfile=/var/run/carbon-relay-1.pid start
stdout_logfile = /var/log/supervisor/%(program_name)s.log
stderr_logfile = /var/log/supervisor/%(program_name)s.log
autorestart = true

[program:carbon-cache-1]
command = /usr/bin/carbon-cache --config=/etc/carbon/carbon.conf --instance=1 --debug --pidfile=/var/run/carbon-cache-1.pid start
stdout_logfile = /var/log/supervisor/%(program_name)s.log
stderr_logfile = /var/log/supervisor/%(program_name)s.log
autorestart = true

[program:carbon-cache-2]
command = /usr/bin/carbon-cache --config=/etc/carbon/carbon.conf --instance=2 --debug --pidfile=/var/run/carbon-cache-2.pid start
stdout_logfile = /var/log/supervisor/%(program_name)s.log
stderr_logfile = /var/log/supervisor/%(program_name)s.log
autorestart = true

[program:graphite-api]
command = gunicorn -b 0.0.0.0:8888 -w 2 --log-level debug graphite_api.app:app
stdout_logfile = /var/log/supervisor/%(program_name)s.log
stderr_logfile = /var/log/supervisor/%(program_name)s.log
autorestart = true